# cups/rootfs/etc/services.d/cups/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
declare exit_code=$1

bashio::log.info "CUPS service stopped with exit code ${exit_code}"