# cups/rootfs/etc/services.d/avahi/run
#!/command/with-contenv bashio
# rootfs/etc/services.d/avahi/run
# Run Avahi daemon under s6 supervision

set -e

bashio::log.info "Starting Avahi daemon..."

# Ensure Avahi runtime dir exists
mkdir -p /var/run/avahi-daemon

# Start Avahi in foreground (Alpine build supports --debug, not --no-daemon)
exec avahi-daemon --no-drop-root --no-chroot
